# systemd/.config/systemd/user/dconf-backup.service
[Unit]
Description=Backup dconf to text file

[Service]
ExecStart=%h/.local/bin/dconf-backup.sh
