# systemd/.config/systemd/user/rss2email.timer
[Unit]
Description=Run rss2email daily 

[Timer]
OnCalendar=01:23
Persistent=true

[Install]
WantedBy=timers.target
