# systemd/.config/systemd/user/rss2email.service
[Unit]
Description=Run rss2email

[Service]
Type=oneshot
ExecStart=%h/.local/bin/rss2email
